# tests/monitoring/prune-synthetic.service
[Unit]
Description=Prune synthetic-monitor rows from Postgres (SYNMON-prefixed VINs older than 24h)
After=network.target

[Service]
Type=oneshot
ExecStart=%h/insurance-app/tests/monitoring/prune-synthetic.sh
StandardOutput=journal
StandardError=journal
